# services/agent0.service
[Unit]
Description=Agent-0 AutoGen Council
After=network-online.target nvidia-persistenced.service
Wants=network-online.target
Requires=nvidia-persistenced.service

[Service]
Type=simple
User=agent0
Group=agent0
Environment=SWARM_GPU_PROFILE=rtx4070
Environment=PYTORCH_CUDA_ALLOC_CONF=max_split_size_mb:128
Environment=PYTHONPATH=/opt/agent0
Environment=AGENT0_SERVICE_MANAGED=true
WorkingDirectory=/opt/agent0
ExecStartPre=/usr/bin/python models/prefetch.py
ExecStart=/usr/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --log-level info
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
LimitNOFILE=65536

# Resource limits
MemoryMax=16G
CPUQuota=200%

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/agent0/logs /opt/agent0/data

[Install]
WantedBy=multi-user.target 