# services/agent0.service
[Unit]
Description=Agent-0 AutoGen Council
After=network.target
Wants=network.target

[Service]
Type=simple
User=agent0
Group=agent0
Environment=SWARM_GPU_PROFILE=rtx4070
Environment=PYTHONPATH=/opt/agent0
WorkingDirectory=/opt/agent0
ExecStart=/usr/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --log-level info
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
LimitNOFILE=65536

# Resource limits
MemoryMax=16G
CPUQuota=200%

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/agent0/logs /opt/agent0/data

[Install]
WantedBy=multi-user.target 